# blocklist.timer
[Unit]
Description=Refresh blocklists routinely

[Timer]
OnUnitActiveSec=2d

[Install]
WantedBy=multi-user.target

