# blocklist.timer
[Unit]
Description=Refresh blocklists routinely

[Timer]
OnUnitActiveSec=2d

[Install]
Wanted-By=multi-user.target

